# templates/kurento/kurento-media-server.target
[Unit]
Description=Kurento Media Server daemon

[Install]
WantedBy=multi-user.target
